# setup/wifry-frontend.service
[Unit]
Description=WiFry Frontend Web Server
After=wifry-backend.service
Wants=wifry-backend.service

[Service]
Type=simple
User=wifry
Group=wifry
WorkingDirectory=/opt/wifry/frontend

# Serve the built React app using Python's http.server
# Proxying to backend is handled by the backend CORS config + client-side API calls
ExecStart=/usr/bin/python3 -m http.server 3000 --directory /opt/wifry/frontend/dist
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifry-frontend

[Install]
WantedBy=multi-user.target
